# HW8-bash/files/nla.timer
[Unit]
Description=Run NginxLogAnalyzer hourly

[Timer]
OnCalendar=hourly
Persistent=true     
 
[Install]
WantedBy=timers.target
